# www-servers/mathopd/files/init
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting mathopd"
	/sbin/start-stop-daemon --start \
	        --exec /usr/bin/mathopd --pidfile /var/run/mathopd.pid \
		-- -f ${MATHOPD_CONF}
	eend $?
}

stop() {
	ebegin "Stopping mathopd"
	start-stop-daemon --stop --quiet --exec /usr/bin/mathopd
	rm -f /var/run/mathopd.pid
	eend $?
}

